# crontab
DISPLAY=:0.0
00 * * * * notify-send "時報" "$(date +\%H):00!"
